# refresh-roothints.service
[Unit]
Description=Refresh PhantomDNS Root Hints

[Service]
Type=oneshot
ExecStart=/etc/pdns/refresh-roothints.sh
